# systemd/acpi-wakeup-persistd.path
[Unit]
Description=Watch /proc/acpi/wakeup for changes

[Path]
PathChanged=/proc/acpi/wakeup

[Install]
WantedBy=paths.target
